# deploy/nullchat-restart.service
[Unit]
Description=Restart nullchat to clear in-memory state

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart nullchat.service
